# TK1_Sample_File_System/etc/init/ssh.conf
# ssh - OpenBSD Secure Shell server
#
# The OpenSSH server provides secure shell access to the system.

description	"OpenSSH server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

env SSH_SIGSTOP=1
expect stop

# 'sshd -D' leaks stderr and confuses things in conjunction with 'console log'
console none

pre-start script
    test -x /usr/sbin/sshd || { stop; exit 0; }
    test -e /etc/ssh/sshd_not_to_be_run && { stop; exit 0; }
    test -f /etc/ssh/ssh_host_dsa_key || { echo "Generating DSA key..." ; ssh-keygen -q -t dsa -N '' -f /etc/ssh/ssh_host_dsa_key; }
    test -f /etc/ssh/ssh_host_rsa_key || { echo "Generating RSA key..." ; ssh-keygen -q -t rsa -N '' -f /etc/ssh/ssh_host_rsa_key; }
    test -f /etc/ssh/ssh_host_ecdsa_key || { echo "Generating ECDSA key..." ; ssh-keygen -q -t ecdsa -N '' -f /etc/ssh/ssh_host_ecdsa_key; }
    mkdir -p -m0755 /var/run/sshd
end script

# if you used to set SSHD_OPTS in /etc/default/ssh, you can change the
# 'exec' line here instead
exec /usr/sbin/sshd -D
